# quikd/quikav-daemon.socket.in
[Unit]
Description=Socket for Quik AntiVirus userspace daemon
Documentation=man:quikd(8) man:quikd.conf(5) http://www.quikav.net/lang/en/doc/
# Check for database existence
ConditionPathExistsGlob=@DBDIR@/main.{c[vl]d,inc}
ConditionPathExistsGlob=@DBDIR@/daily.{c[vl]d,inc}

[Socket]
ListenStream=/run/quikav/quikd.ctl
#ListenStream=127.0.0.1:1024
SocketUser=quikav
SocketGroup=quikav
RemoveOnStop=True

[Install]
WantedBy=sockets.target
